# assemblyline/assemblyline/al/install/etc/init/expiry.conf
description "AL Expiry Loader"

start on runlevel [2345]
stop on runlevel [016]


script
    . /etc/default/al
    exec sudo -u $AL_USER \
    AL_DATASTORE=${AL_DATASTORE} AL_ROOT=${AL_ROOT} PYTHONPATH=${PYTHONPATH} \
    ${AL_ROOT}/pkg/assemblyline/al/run/expiry.py
end script
