# assemblyline/assemblyline/al/run/vmbootstrap/hostagent-withbootstrap.conf
description "AL HostAgent"

start on net-device-up IFACE!=lo
stop on runlevel [016]


# Use the SIGINT signal first instead of SIGTERM to give agent 
# the chance to cleanly come down.
kill signal SIGINT
kill timeout 15

pre-start script
    . /etc/default/al
    /opt/al/pkg/assemblyline/al/run/vmbootstrap/hostagent-bootstrap-stage-0.py
    /opt/al/pkg/assemblyline/al/run/vmbootstrap/hostagent-bootstrap-stage-1.py

    . /etc/default/al
    sudo -u ${AL_USER} \
     PYTHONPATH=${PYTHONPATH} AL_DATASTORE=${AL_DATASTORE} \
     /opt/al/pkg/assemblyline/al/run/hostagent.py --sysprep
end script

respawn

script
  . /etc/default/al

  exec sudo -u ${AL_USER} \
     PYTHONPATH=${PYTHONPATH} AL_DATASTORE=${AL_DATASTORE} \
     /opt/al/pkg/assemblyline/al/run/hostagent.py
end script
